# linux_env/profile.d/python.sh
PYTHONSTARTUP=~/bin/tools/python_startup.py
export PYTHONSTARTUP
