# linux_env/profile.d/python.sh
PYTHONSTARTUP=~/bin/python_startup.py
export PYTHONSTARTUP
